# src/tiny-ramdisk-ramfs.service
[Unit]
Description=Tiny RAMDisk

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/sleep 5
ExecStart=%h/.local/bin/tiny-ramdisk-ramfs-startup.sh
ExecStop=%h/.local/bin/tiny-ramdisk-ramfs-shutdown.sh

[Install]
WantedBy=multi-user.target
